# infrastructure/supervisor/horizon.conf
[program:horizon]
command=php /var/www/html/apps/api/artisan horizon
directory=/var/www/html/apps/api
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
stdout_logfile=/var/log/supervisor/horizon.log
stderr_logfile=/var/log/supervisor/horizon-error.log
